# linux/systemd/unit_timer_test/unit_timer_test.timer
[Unit]
Description=unit_timer_test

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target
