# linux/systemd/unit_timer_test/unit_timer_test.service
[Unit]
Description=unit_timer_test

[Service]
ExecStart=/home/cytopyge/_test/unit_timer_test.sh
